# script/server/0_local_services/sysom-redis.ini
[program:sysom-redis]
directory = /usr/local/sysom/server/redis
command=/usr/local/sysom/server/redis/redis-server /usr/local/sysom/server/redis/redis.conf
startsecs=3
autostart=true
autorestart=true
stderr_logfile=/usr/local/sysom/server/logs/sysom-redis-error.log
stdout_logfile=/usr/local/sysom/server/logs/sysom-redis.log
