# tools/deploy/channel-service.ini
[program:channel-service]
directory = /usr/local/sysom/server/target/sysom_server/sysom_channel
command=/usr/local/sysom/server/virtualenv/bin/gunicorn -c ./conf/channel_gunicorn.py sysom_channel.wsgi:application
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/supervisord-channel-error.log
stdout_logfile=/usr/local/sysom/server/logs/supervisord-channel.log
