# tools/deploy/sysom-vmcore.ini
[program:sysom-vmcore]
directory = /usr/local/sysom/server/target/sysom_server/sysom_vmcore
command=/usr/local/sysom/server/virtualenv/bin/gunicorn -c ./conf/vmcore_gunicorn.py sysom_vmcore.wsgi:application
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/sysom-vmcore-error.log
stdout_logfile=/usr/local/sysom/server/logs/sysom-vmcore.log
