# tools/deploy/sysom-diagnosis.ini
[program:sysom-diagnosis]
directory = /usr/local/sysom/server/target/sysom_server/sysom_diagnosis
command=/usr/local/sysom/server/virtualenv/bin/gunicorn -c ./conf/diagnosis_gunicorn.py sysom_diagnosis.wsgi:application
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/supervisord-diagnosis-error.log
stdout_logfile=/usr/local/sysom/server/logs/supervisord-diagnosis.log
