# tools/deploy/sysom.ini
[program:sysom]
directory = /home/sysom/target/sysom_api
command=/home/sysom/virtualenv/bin/uwsgi --ini /home/sysom/target/sysom_api/uwsgi.ini --catch-exceptions --protocol=http
autostart=true
autorestart=true
redirect_stderr=true
stopasgroup=true
stdout_logfile=/home/sysom/logs/supervisor.log
environment=PATH="/home/sysom/virtualenv/bin"
