# tools/deploy/sysom-channel.ini
[program:sysom-channel]
directory = /usr/local/sysom/server/target/sysom_server/sysom_channel_v2
command=/usr/local/sysom/server/virtualenv/bin/uvicorn main:app --port 7003
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/supervisord-channel-error.log
stdout_logfile=/usr/local/sysom/server/logs/supervisord-channel.log
