# tools/deploy/task-service.ini
[program:task-service]
directory = /usr/local/sysom/server/target/sysom_api
command=/usr/local/sysom/server/virtualenv/bin/gunicorn -c ./conf/task_gunicorn.py sysom.wsgi:application
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/supervisord-task-error.log
stdout_logfile=/usr/local/sysom/server/logs/supervisord-task.log