# tools/deploy/task-executor-service.ini
[program:task-executor-service]
directory = /usr/local/sysom/server/target/sysom_api
command=/usr/local/sysom/server/virtualenv/bin/python -m microservices.task_executor.app
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/supervisord-task-executor-error.log
stdout_logfile=/usr/local/sysom/server/logs/supervisord-task-executor.log
