# script/server/6_sysom_hotfix_builder/sysom-hotfix-builder.ini
[program:sysom-hotfix-builder]
directory = /usr/local/sysom/server/target/sysom_server/sysom_hotfix_builder
command=/usr/local/sysom/server/virtualenv/bin/python3 builder.py
startsecs=3
autostart=true
autorestart=true
environment=PATH=/usr/local/sysom/server/virtualenv/bin:%(ENV_PATH)s
stderr_logfile=/usr/local/sysom/server/logs/sysom-hotfix-builder-error.log
stdout_logfile=/usr/local/sysom/server/logs/sysom-hotfix-builder.log
