# script/server/1_sysom_monitor/sysom-prometheus.ini
[program:sysom-prometheus]
directory = /usr/local/sysom/server/monitor/prometheus
command=/usr/local/sysom/server/monitor/prometheus/prometheus
startsecs=3
autostart=true
autorestart=true
stderr_logfile=/usr/local/sysom/server/logs/sysom-prometheus-error.log
stdout_logfile=/usr/local/sysom/server/logs/sysom-prometheus.log
