# script/server/7_sysom_demo/sysom-demo.ini
[program:sysom-demo]
directory = /usr/local/sysom/server/target/sysom_server/sysom_demo
command=/usr/local/sysom/server/virtualenv/bin/uvicorn main:app --port 7008
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/sysom-demo-error.log
stdout_logfile=/usr/local/sysom/server/logs/sysom-demo.log
