# script/server/1_sysom_monitor_server/sysom-monitor-server.ini
[program:sysom-monitor-server]
directory = /usr/local/sysom/server/target/sysom_server/sysom_monitor_server
command=/usr/local/sysom/server/virtualenv/bin/uvicorn main:app --port 7009
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/sysom-monitor-server-error.log
stdout_logfile=/usr/local/sysom/server/logs/sysom-monitor-server.log
