# script/server/6_sysom_hotfix/sysom-hotfix.ini
[program:sysom-hotfix]
directory = /usr/local/sysom/server/target/sysom_server/sysom_hotfix
command=/usr/local/sysom/server/virtualenv/bin/gunicorn -c ./conf/hotfix_gunicorn.py sysom_hotfix.wsgi:application
startsecs=3
autostart=true
autorestart=true
environment=PATH="/usr/local/sysom/server/virtualenv/bin/"
stderr_logfile=/usr/local/sysom/server/logs/sysom-hotfix-error.log
stdout_logfile=/usr/local/sysom/server/logs/sysom-hotfix.log
